# templates/etc/profile.d/atom.sh.j2
#!/bin/bash
# this file managed by ansible
# AtoM environment variables (required by CLI tools)
# Not adding PATH or null env vars
{% for key, value in atom_pool_php_envs.items() %}
{% if key != "PATH" and value != "" -%}
export {{ key }}="{{ value }}"
{% endif %}
{% endfor %}

{% if ansible_os_family == "RedHat" and php_version|int < 74 %}
# enable SCLs for node and php
source scl_source enable rh-nodejs6 rh-php{{ php_version }}
{% endif %}

