# templates/etc/init/atom-worker.conf
# Ansible managed file, do not edit directly
description "AtoM worker upstart service"

start on startup
stop on shutdown

# Respawn the job up to 10 times within a 60 seconds period. If the job
# exceeds these values, it will be stopped and marked as failed. We delay
# the respawn 3 seconds. This is a workaround to wait for other services
# to start during the boot process. It should be handled by the application.
respawn
respawn limit 10 60
post-stop exec sleep 3

env LOCATION={{ atom_path }}
env LOGFILE={{ atom_path }}/log/worker.log
# AtoM PHP pool vars 
{% for key, value in atom_pool_php_envs.items() %}
{% if value != "" -%}
    env {{ key }}="{{ value }}"
{% endif %}
{% endfor %}


setuid {{ atom_user }}
setgid {{ atom_group }}

script

  php \
    -d memory_limit=-1 \
    -d error_reporting="E_ALL" \
{% if atom_worker_old_config|bool %}
      ${LOCATION}/symfony tools:gearman-worker >> ${LOGFILE} 2>&1
{% else %}
      ${LOCATION}/symfony jobs:worker >> ${LOGFILE} 2>&1
{% endif %}

end script
