# templates/etc/init/drmc-mock.conf
# Ansible managed file, do not edit directly
description "DRMC worker upstart service"

start on startup
stop on shutdown

# Respawn the job up to 10 times within a 60 seconds period. If the job
# exceeds these values, it will be stopped and marked as failed. We delay
# the respawn 3 seconds. This is a workaround to wait for other services
# to start during the boot process. It should be handled by the application.
respawn
respawn limit 10 60
post-stop exec sleep 3

{% if atom_revision_directory|bool %}
env LOCATION={{ atom_path }}/{{ atom_revision_directory_latest_symlink_dir }}/plugins/arDrmcPlugin/frontend/mock_api
{% else %}
env LOCATION={{ atom_path}}/plugins/arDrmcPlugin/frontend/mock_api
{% endif %}

setuid {{ atom_user }}
setgid {{ atom_user }}

script

  cd $LOCATION && node app.js

end script
